# ip-blocked.service
[Unit]
Description=IP-block service

[Service]
ExecStart=/<location>/blocked.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
